# pobDNSupdate.timer
[Unit]
Description=A dynamic dns service on your web space (Timer)

[Timer]
OnUnitActiveSec=300s
OnBootSec=10s

[Install]
WantedBy=timers.target
